# scripts/wg-status-sync.service
[Unit]
Description=WireGuard VPN Status Sync
After=network.target docker.service

[Service]
Type=oneshot
ExecStart=/opt/anvil/scripts/wg-status-sync.sh
